# deployment/cookbooks/dmon/templates/default/start_component.service.erb
[Unit]
Description=<%=@componentName%> service
[Service]
ExecStart=<%=@home%>/bin/start.sh
ExecStop=<%=@home%>/bin/stop.sh
WorkingDirectory=<%=@home%>/bin
User=root
Type=forking
[Install]
WantedBy=multi-user.target
